# deploy/sosenki.service.template
# SOSenki systemd service template
# Generated by: make install-service
# Location: /etc/systemd/system/sosenki.service

[Unit]
Description=SOSenki Telegram Bot + Mini App
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=${USER}
Group=${USER}
WorkingDirectory=${INSTALL_DIR}
EnvironmentFile=${INSTALL_DIR}/.env
ExecStart=${INSTALL_DIR}/.venv/bin/python -m src.main --mode webhook
Restart=on-failure
RestartSec=5
StandardOutput=append:${INSTALL_DIR}/logs/server.log
StandardError=append:${INSTALL_DIR}/logs/server.log

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=${INSTALL_DIR}/logs ${INSTALL_DIR}/backups ${INSTALL_DIR}

[Install]
WantedBy=multi-user.target
